# opensuse/packages/tf300/lib/systemd/system/tf300-cpu.service
[Unit]
Description=Set Tegra parameters
DefaultDependencies=no
After=sysinit.target
Before=shutdown.target
Conflicts=shutdown.target
TimeoutStartSec=20

[Service]
Type=oneshot
ExecStart=/opt/tf300/cpu.sh
StandardOutput=syslog

[Install]
WantedBy=multi-user.target
